# orchestrator.service
[Unit]
Description=Service Orchestrator
[Service]
ExecStart=python3 orchestrator.py
Restart=always
[Install]
WantedBy=multi-user.target
